# deploy/firewall-manager.service
[Unit]
Description=Firewall Manager
After=network-online.target
Wants=network-online.target

[Service]
User=firewall-manager
Group=firewall-manager
EnvironmentFile=/etc/firewall-manager/env
ExecStart=/usr/local/bin/firewall-manager
Restart=on-failure
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run /etc/ufw /var/lib/firewall-manager /var/log/firewall-manager

[Install]
WantedBy=multi-user.target
